# var/systemd/infinitysky-cron.timer
[Unit]
Description=Runs cron job every minute

[Timer]
OnCalendar=*-*-* *:*:00
Unit=infinitysky-cron.service

[Install]
WantedBy=timers.target
